# docker/supervisor/rabbitmq.conf
[program:rabbitmq]
command=/usr/sbin/rabbitmq-server
numprocs=1
autostart=true
autorestart=true
stdout_logfile=/var/log/rabbitmq/rabbitmq.log
stderr_logfile=/var/log/rabbitmq/rabbitmq.log